# mercury/var/spool/cron/crontabs/mble
@daily pihole updateGravity
@daily pihole -up
@daily apt update && apt full-upgrade
